# vncsdk-1.7.0/samples/serviceServerLinux/upstart/vncsdkServer.conf
description   "VNC SDK Linux service service sample"


start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/local/bin/serviceServerLinux
